# kafka/supervisor/01_kafka.conf
[program:kafka]
command=/usr/local/sbin/start-kafka.sh
autostart=true
autorestart=true
