# full/supervisor/Social-Finance.conf
[program:Social-Finance]
autostart=true
command=mvn jetty:deploy-war -Djetty.port=8081
directory=/opt/OBP/Social-Finance
killasgroup=true
startsecs=5
stderr_logfile = /var/log/supervisor/Social-Finance-stderr.log
stdout_logfile = /var/log/supervisor/Social-Finance-stdout.log
stopasgroup=true
